# systemd/sound-led-daemon.service
[Unit]
Description=Sound Machine LED Control Daemon
After=network.target
Before=sound-trigger.service

[Service]
Type=simple
User=soundconsole
WorkingDirectory=/home/soundconsole/sound-machine
ExecStart=/usr/bin/python3 /home/soundconsole/sound-machine/daemon/led_daemon.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

