# scripts/mss-health.service
[Unit]
Description=MSS Health Check
After=docker.service

[Service]
Type=oneshot
ExecStart=/home/shschubert/Projects/sound-machine/scripts/healthcheck.sh
StandardOutput=journal
StandardError=journal
