# scripts/mss-health.timer
[Unit]
Description=Run MSS Health Check every 5 minutes

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min
AccuracySec=30s

[Install]
WantedBy=timers.target
